# systemd/dental-celery-beat.service
[Unit]
Description=Odonta Index AI Celery Beat Scheduler
After=network.target redis.service dental-backend.service

[Service]
Type=simple
User=deploy
WorkingDirectory=/opt/dental-ai
ExecStart=/opt/dental-ai/.venv/bin/celery -A app.celery_app beat --loglevel=info
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
